# crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
17  7   */2 *   * /app/run.sh
5   5   *   */2 * apk add --update-cache --upgrade apk-tools ca-certificates
